# resolver/pkg/pzresolver.service
[Unit]
Description=Promtuz resolver (relay discovery service)
Documentation=https://github.com/promtuz/promtuz
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pzresolver
Group=pzresolver
ExecStart=/usr/bin/pzresolver --config /etc/promtuz/resolver.toml
Restart=on-failure
RestartSec=2

# Least-privilege hardening. /usr stays read-only; /etc/promtuz is writable so
# the resolver can generate its key, write a CSR, and read the enrolled cert.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_INET AF_INET6
ConfigurationDirectory=promtuz
ReadWritePaths=/etc/promtuz

[Install]
WantedBy=multi-user.target
