# resolver/pkg/scripts/postinst
#!/bin/sh
set -e

if ! getent passwd pzresolver >/dev/null 2>&1; then
    adduser --system --group --no-create-home --home /nonexistent pzresolver
fi

install -d -o root       -g root       -m 0755 /etc/promtuz
install -d -o root       -g root       -m 0711 /etc/promtuz/keys
install -d -o root       -g root       -m 0755 /etc/promtuz/certs
install -d -o pzresolver -g pzresolver -m 0700 /etc/promtuz/keys/resolver

systemctl daemon-reload >/dev/null 2>&1 || true

# Fresh install ($2 empty): enable + start now. The resolver generates its key,
# writes a CSR, and waits for the signed cert (ensure_enrolled never
# crash-loops), so starting before enrollment is safe and gets the CSR written
# immediately. Upgrade: try-restart to pick up the new binary if running.
if [ "$1" = "configure" ] && [ -z "$2" ]; then
    systemctl enable --now pzresolver >/dev/null 2>&1 || true
else
    systemctl try-restart pzresolver >/dev/null 2>&1 || true
fi

cat <<'EOF'
pzresolver installed, started, and waiting for its cert. One step left, on your CA box:
  1. certgen sign /etc/promtuz/keys/resolver/resolver.csr
  2. drop the signed cert at /etc/promtuz/certs/resolver.crt   # starts serving automatically
EOF

exit 0